# installer/macos/templates/com.agentrooms.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentrooms.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/usr/local/lib/node_modules/agentrooms/dist/cli/node.js</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8080</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/__USER__/Library/Logs/agentrooms-error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/__USER__/Library/Logs/agentrooms-output.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/__USER__</string>
</dict>
</plist>